# tedge/etc/init.d/tedge-mgm
#!/sbin/openrc-run

name=$RC_SVCNAME
description="tedge-mgm runs web ui to manage the thin edge."

command="/app/tedge/server/start.sh"
command_args=""
command_user="tedge-mgm"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
output_log="/var/log/$RC_SVCNAME/$RC_SVCNAME.log"
error_log="/var/log/$RC_SVCNAME/$RC_SVCNAME.err"

depend() {
    need net
}

start_pre() {
    /bin/touch "/var/log/$RC_SVCNAME/$RC_SVCNAME.log"
    /bin/chmod 744 "/var/log/$RC_SVCNAME/$RC_SVCNAME.log"
    /bin/touch "/var/log/$RC_SVCNAME/$RC_SVCNAME.err"
    /bin/mkdir -p /var/log
    /bin/mkdir -m 744 "/var/log/$RC_SVCNAME"
    /bin/chown tedge-mgm:tedge-mgm "/var/log/$RC_SVCNAME"
}
