# tedge/etc/init.d/tedge-mgm
#!/sbin/openrc-run

name=$RC_SVCNAME
description="ui to manage the thin-edge.io locally"

command="/usr/bin/tedge-mgm"
command_args=""
command_user="tedge-mgm"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
output_log="/var/log/$RC_SVCNAME/$RC_SVCNAME.log"
error_log="/var/log/$RC_SVCNAME/$RC_SVCNAME.err"

depend() {
    need net
}

start_pre() {
    # /bin/mkdir -p /var/log
    /bin/mkdir -m 755 -p "/var/log/$RC_SVCNAME"
    /bin/touch "/var/log/$RC_SVCNAME/$RC_SVCNAME.log"
    /bin/chmod 755 "/var/log/$RC_SVCNAME/$RC_SVCNAME.log"
    /bin/touch "/var/log/$RC_SVCNAME/$RC_SVCNAME.err"
    /bin/chmod 755 "/var/log/$RC_SVCNAME/$RC_SVCNAME.err"
    /bin/chown tedge-mgm:tedge-mgm -R "/var/log/$RC_SVCNAME"
}
